# systemd/user/obsidian-stats-helper.service
[Unit]
Description=Harvest Obsidian plugin metadata
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/src/obsidian-stats-helper
EnvironmentFile=%h/.config/obsidian-stats-helper/env
Environment=PATH=%h/.bun/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/src/obsidian-stats-helper/scripts/harvest-pi.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
