# systemd/user/obsidian-stats-helper.timer
[Unit]
Description=Run Obsidian plugin metadata harvest

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
